# src/templates/systemd/healthcheck.service.tmpl
[Unit]
Description=RayLink ${NODE_ROLE} Reality health check
After=network-online.target ${XRAY_SERVICE}
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-${HEALTHCHECK_ENV_FILE}
ExecStart=${HEALTHCHECK_EXEC}
